# debian/nazs.upstart
start on runlevel [2345]
stop on runlevel [!2345]

script
  /usr/bin/nazs runserver --nothreading --noreload 0.0.0.0:8000
end script
